# doc/qresp.socket
# /etc/systemd/system/qresp.service

[Unit]
Description = qresp socket

[Socket]
ListenStream = /run/qresp.sock

[Install]
WantedBy = sockets.target
